# root/etc/services.d/openssh-server/run
#!/usr/bin/with-devcontenv bash

exec /usr/sbin/sshd -D
