# system/devcont/dist/etc/profile.d/00-devcont-profile.sh
#!/bin/bash

. /usr/local/sbin/devcont-import-all
