# packages/docker/dist/etc/services.d/docker/run
#!/usr/local/sbin/with-devcontenv bash

exec dockerd
