# root/etc/services.d/docker/run
#!/usr/bin/with-devcontenv bash

exec dockerd 2>&1 | sudo -u ${USER_NAME} tee /mnt/data/docker/docker.log > /dev/null
